# puppet/modules/hhvm/templates/hhvm.conf.erb
# hhvm - HipHop Virtual Machine
# Runs HHVM in server mode.
description "HipHop Virtual Machine for PHP"

start on (local-filesystems and net-device-up IFACE!=lo)

setuid www-data
setgid www-data

exec /usr/bin/hhvm --config /etc/hhvm/server.hdf --mode server

respawn
# vim: set ft=upstart:
