# puppet/modules/hhvm/templates/hhvm.conf.erb
# hhvm - HipHop Virtual Machine
# Runs HHVM in server mode.
description "HipHop Virtual Machine for PHP"

start on (local-filesystems and net-device-up IFACE!=lo)

pre-start script
  mkdir -p -m0755 /var/run/hhvm
  chown -R www-data:www-data /var/run/hhvm
  mkdir -p -m0755 <%= @logroot %>/hhvm
  # chowns don't work on Vagrant's NFS mounted logroot, but they also should
  # not be needed there due to the uid mapping configuration
  chown -R www-data:www-data <%= @logroot %>/hhvm || /bin/true

  # Update the symlink `/usr/lib/x86_64-linux-gnu/hhvm/extensions/current` so
  # it points to the directory whose name matches the extension API version
  # (e.g. `/usr/lib/x86_64-linux-gnu/hhvm/extensions/20140829`).
  API_VERSION="$(/usr/bin/hhvm --version | grep -Po '(?<=API: ).+')"
  ( cd /usr/lib/x86_64-linux-gnu/hhvm/extensions; ln -Trsf "$API_VERSION" current; )
end script

env REMOTE_ADDR="127.0.0.1"

exec /sbin/start-stop-daemon --quiet --start \
  --chuid www-data:www-data \
  --pidfile /var/run/hhvm/hhvm.pid \
  --startas /usr/bin/hhvm -- \
    --config /etc/hhvm/config.hdf \
    --config /etc/hhvm/fcgi/php.ini \
    --mode server -v Eval.Jit=true

# Don't limit the size of core dumps.
limit core unlimited unlimited

# Increase the maximum number of open files.
limit nofile 65536 65536

# When 'hhvm.server.graceful_shutdown_wait' is set to a positive
# integer, HHVM will perform a graceful shutdown on SIGHUP.
kill signal HUP

respawn
# vim: set ft=upstart:
