# puppet/modules/statsd/templates/upstart.conf.erb
# vim: set ft=upstart:

# This file is managed by Puppet
# Upstart init script for statsd

description "statsd service"

start on mediawiki-ready
stop on runlevel [016]

setuid "www-data"
setgid "www-data"

env NODE_PATH="<%= @dir %>/node_modules"

chdir "<%= @dir %>"
exec /usr/bin/nodejs stats.js <%= @dir %>/config.js > <%= @logdir %>/statsd.log

respawn

