# puppet/modules/activemq/templates/activemq-upstart.conf.erb
description "ActiveMQ STOMP server"

start on (local-filesystems and net-device-up IFACE!=lo)

setuid activemq
setgid activemq

exec <%= @pkgdir %>/bin/activemq-admin
