# puppet/modules/eventlogging/templates/devserver.upstart.erb
#####################################################################
### THIS FILE IS MANAGED BY PUPPET
#####################################################################

description "EventLogging Dev Server"

start on web
respawn

setuid www-data
setgid www-data

exec <%= @eventlogging_path %>/virtualenv/bin/eventlogging-devserver --append-to <%= @output_file %>
