# puppet/modules/eventlogging/templates/service.upstart.erb
#####################################################################
### THIS FILE IS MANAGED BY PUPPET
#####################################################################

description "EventLogging Produce HTTP Service"

start on web
respawn

setuid www-data
setgid www-data

exec <%= @eventlogging_path %>/virtualenv/bin/eventlogging-service @<%= @config_file %>
