# puppet/modules/zotero/templates/systemd.service.erb
# NOTE: This file is managed by Puppet
# Systemd unit for Zotero
[Unit]
Description="Zotero service"
After=network.target

[Service]
# up ulimit -n a bit
LimitNOFILE=10000
User=www-data
Group=www-data
Environment="LD_LIBRARY_PATH=/usr/lib/xulrunner-24.0/"
Restart=always
RestartSec=2s
# wait 60 seconds for a graceful restart before killing the master
TimeoutStopSec=60
# The xulrunner-dev package installs a non-executable
# /usr/lib/xulrunner-devel-24.0/sdk/bin/xpcshell
# Work around that by invoking ld-linux-x86-64.so.2 first.
ExecStart=/lib64/ld-linux-x86-64.so.2 /usr/lib/xulrunner-devel-24.0/sdk/bin/xpcshell -g /usr/lib/xulrunner-24.0/ -a <%= @base_path %>/translation-server -mn -v 180 <%= @base_path %>/translation-server/translation-server/init.js >> <%= @log_file %> 2>&1

[Install]
WantedBy=multi-user.target
