# puppet/modules/restbase/templates/upstart.conf.erb
# vim: set ft=upstart:

# Upstart job configuration for RESTBase
# This file is managed by Puppet

description "RESTBase REST API service"

start on mediawiki-ready
stop on runlevel [!2345]

setuid "www-data"
setgid "www-data"

env VCAP_APP_PORT="<%= @port %>"
env NODE_PATH="<%= @dir %>/node_modules"

chdir "<%= @dir %>"
exec /usr/bin/nodejs server.js

respawn

